# root/etc/init.d/network_scripts
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1
PROG=/usr/bin/network_switch.sh

start_service() {
    # 0点断网和6点联网任务还是得用Crontab，这个是硬性时间规定
    sed -i '/校园网切换/d' /etc/crontabs/root 2>/dev/null
    echo "0 0 * * * sleep 5 && uci set network.wan.disabled='1' && uci commit network && ifdown wan && logger '校园网切换: 0点已禁用WAN'" >> /etc/crontabs/root
    echo "0 6 * * * uci set network.wan.disabled='0' && uci commit network && ifup wan && logger '校园网切换: 6点已启用WAN'" >> /etc/crontabs/root
    /etc/init.d/cron restart

    procd_open_instance
    procd_set_param command "$PROG" "monitor"
    procd_set_param stdout 1
    procd_set_param stderr 1
    # 崩溃后自动重启，间隔 5 秒
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
    procd_close_instance
}

stop_service() {
    # 移除 Crontab 任务
    sed -i '/校园网切换/d' /etc/crontabs/root 2>/dev/null
    /etc/init.d/cron restart
}